# udev/99-hdmi-dsi.rules
ACTION=="change", SUBSYSTEM=="drm", RUN+="/bin/systemctl start hdmi-dsi-fix.service"
